# deployment/macos/com.maktabaa.tracking-rpa.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.maktabaa.tracking-rpa</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd "__REPOSITORY_PATH__" &amp;&amp; exec ./scripts/start-tracking-rpa-macos.command</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/maktabaa-tracking-rpa.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/maktabaa-tracking-rpa-error.log</string>
</dict>
</plist>
